# deploy/taskbrain-finance-api.service
[Unit]
Description=TaskBrain Finance API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/taskbrain/finance/backend
EnvironmentFile=/opt/taskbrain/finance/backend/.env
ExecStart=/opt/taskbrain/finance/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8090
Restart=on-failure
RestartSec=5
User=taskbrain
Group=taskbrain

[Install]
WantedBy=multi-user.target

